# systemd/gaming-vm.service
[Service]
ExecStart=/home/sam/repos/qemupin/vm.sh
Slice=qemu.slice
